# overlays/rootfs/etc/systemd/system/serial-getty@ttyAMA0.service.d/autologin.conf
# NanoOS: autologin root on serial console (QEMU)
# This avoids needing a login prompt for the mobile OS

[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin root --noclear %I $TERM
